# scripts/deploy/autoppia-miner.service
[Unit]
Description=Autoppia Bittensor Miner
After=network.target autoppia-api.service
Wants=network-online.target
Requires=autoppia-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/autoppia-miner
Environment="PATH=/opt/autoppia-miner/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/autoppia-miner/.env
ExecStart=/opt/autoppia-miner/venv/bin/python3 -m miner.miner --wallet.name ${WALLET_NAME} --wallet.hotkey ${WALLET_HOTKEY} --netuid ${SUBNET_UID} --subtensor.network ${NETWORK} --axon.port ${AXON_PORT}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autoppia-miner

# Resource limits
LimitNOFILE=65536
MemoryMax=4G

[Install]
WantedBy=multi-user.target
